# root/etc/services.d/prowlarr/run
#!/usr/bin/with-contenv bash

echo "Starting Prowlarr."
redirect_cmd exec s6-setuidgid abc \
	/app/prowlarr/bin/Prowlarr \
	-nobrowser \
	-data=/config
